# airootfs/etc/systemd/system/home-prepare.service
[Unit]
Description=Prepare home directory
Before=gdm.service

[Service]
ExecStart=/bin/home-prepare

[Install]
WantedBy=multi-user.target
