# airootfs/etc/systemd/system/session-proxy.service
[Unit]
Description=Live Session Proxy

[Service]
ExecStart=/bin/session-proxy

[Install]
WantedBy=multi-user.target
